# packages/clipthread-core/clipthread_core-0.0.1.tar.gz/clipthread_core-0.0.1/copythread.service
[Unit]
Description=CopyThread
After=graphical.target

[Service]
Type=simple
User=lcances
Group=lcances
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/lcances/.Xauthority
ExecStart=/home/lcances/Documents/Perso/python/clipthread/run_table_ui.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target